# palm/qemux86_upstart-pulseaudio.conf
#  This file is part of PulseAudio.
#  Copyright (c) 2013-2018 LG Electronics, Inc.
#  All rights reserved.
#
#  PulseAudio is free software; you can redistribute it and/or modify
#  it under the terms of the GNU Lesser General Public License as published
#  by the Free Software Foundation; either version 2.1 of the License,
#  or (at your option) any later version.
#
#  PulseAudio is distributed in the hope that it will be useful, but
#  WITHOUT ANY WARRANTY; without even the implied warranty of
#  MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the GNU
#  General Public License for more details.
#
#  You should have received a copy of the GNU Lesser General Public License
#  along with PulseAudio; if not, see <http://www.gnu.org/licenses/>.


description "pulseaudio server"

start on init-boot-done
stop on started start_update

respawn

# Comment this line out to suppress logs on the console
#console output

# set hardware volumes when running in emulator
pre-start script
        rm  -rf ~/.pulse/*
#group and permissions  changing for alsa device
# should be removed after  fixing GF-52390
        chgrp -R audio /dev/snd/*
        chmod -R 660 /dev/snd/*

end script

exec ionice -c1 /usr/bin/pulseaudio --log-target=syslog --resample-method=palm
